# files/nscd/nscd.ini
[program:nscd]
command=/usr/sbin/nscd -F
priority=1
directory=/
autostart=true
autorestart=true

[program:nscd_invalidate]
command=/usr/local/bin/nscdclear.sh
priority=2
directory=/
autostart=true
autorestart=true
