# files/dnsmasq/dnsmasq.ini
[program:dnsmasq]
command=/usr/sbin/dnsmasq -k
priority=1
directory=/
autostrart=true
autorestart=true